# kernelctf/server/kernelctf.service
[Unit]
Description=kernelCTF

[Service]
ExecStart=/home/poprdi/service.sh
User=kernelctf
Group=kernelctf

[Install]
WantedBy=multi-user.target